# etc/systemd/system/salyut-bbs-web.service
[Unit]
Description=salyut.one BBS read-only web viewer
After=network.target salyut-bbsd.service
Requires=salyut-bbsd.service

[Service]
Type=simple
User=salyut-web
Group=salyut-bbs-web
ExecStart=/usr/local/bin/salyut-bbs-web
Restart=on-failure
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
